# kvirt/cluster/openshift/contrail-ca-patch.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/ca-patch.sh

[Install]
WantedBy=multi-user.target
